# supervisor/mysql-master.ini
[program:mysql-master]
command = scs-docker --log-level silly run ../manifest-vagrant.yaml manifest.yaml
directory = /vagrant/mysql-master
autorestart = true
