# supervisor/mysqlslave.ini
[program:mysqlslave]
command = $PWD/bin/handle-service run mysqlslave
directory = $PWD
autorestart = true
autostart = false
