# supervisor/mysqlmaster.ini
[program:mysqlmaster]
command = $PWD/bin/handle-service run mysqlmaster
directory = $PWD
autorestart = true
autostart = false
