# supervisor/mysql-slave.ini
[program:mysql-slave]
command = scs-docker --log-level silly run ../manifest-vagrant.yaml manifest.yaml
directory = /vagrant/mysql-slave
