# supervisor/scs-disco-server.ini
[program:scs-disco-server]
command = scs-disco-server --log-level silly 0.0.0.0
autorestart = true
autostart = false
