# supervisor/wordpress.ini
[program:wordpress]
command = $PWD/bin/handle-service run wordpress
directory = $PWD
autorestart = true
autostart = false
