# supervisor/wordpress.ini
[program:wordpress]
command = scs-docker --log-level silly run ../manifest-vagrant.yaml manifest.yaml
directory = /vagrant/wordpress
autorestart = true
